# .config/systemd/user/midiangel.service
[Unit]
Description=sdrangel midi interface

[Service]
ExecStart=/home/myon/shack/sdrangel/midiangel
