# .config/systemd/user/notify-send@.service
[Service]
Type=oneshot
ExecStart=/usr/bin/notify-send "systemd user service" %i
