# .config/systemd/user/wsjtx-log.service
[Service]
Type=oneshot
WorkingDirectory=/home/myon/log/wsjtx
ExecStart=/home/myon/log/wsjtx/import

[Unit]
OnFailure=notify-send@wsjtx-log.service
