# .config/systemd/user/wsjtx-log.service
[Service]
Type=oneshot
WorkingDirectory=/home/myon/log/wsjtx
ExecStart=/home/myon/log/wsjtx/import
