# .config/systemd/user/wsjtx-log.path
[Path]
PathChanged=/home/myon/log/wsjtx/wsjtx.log

[Install]
WantedBy=default.target
