# .config/systemd/user/xroot.timer
[Timer]
OnCalendar=hourly

[Install]
WantedBy=default.target
