# .config/systemd/user/xroot.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c ~/bin/xroot

[Install]
WantedBy=default.target
