# .config/systemd/user/xroot.timer
[Timer]
OnCalendar=*:00/20:00

[Install]
WantedBy=default.target
