# .config/systemd/user/xscreensaver.service
[Service]
Type=exec
ExecStart=/usr/bin/xscreensaver -nosplash

[Install]
WantedBy=default.target
